# supervisor.conf
[program:mysql-server]
command = /usr/bin/mysqld_safe
